# etc/systemd/system/rpi-check-restart.service
[Unit]
Description=Restart rpi-check.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart rpi-check.service